# ops/deploy/mchat-backend.service
[Unit]
Description=mchat Backend API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/xiaoxiao/mchat/src/backend
EnvironmentFile=/opt/xiaoxiao/mchat/.env
ExecStart=/opt/xiaoxiao/mchat/src/backend/venv/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 3001
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
